# mqtt-status-updater.service
[Unit]
Description=MQTT Desktop Status Updater

[Service]
ExecStart=/usr/bin/mqtt-status-updater

[Install]
WantedBy=basic.target
